# deploy/cc-im.service
[Unit]
Description=CC-IM Telegram Bot for Claude Code
Documentation=https://github.com/zhaofinger/cc-im
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=%USER%
Group=%GROUP%

# Working directory
WorkingDirectory=%WORKDIR%

# Environment
Environment="PATH=%SERVICE_PATH%"
Environment="HOME=%HOME%"

# Load env file if exists (systemd will fail silently if file doesn't exist due to the - prefix)
EnvironmentFile=-%WORKDIR%/.env

# Command to run
ExecStart=%BUN_PATH% run src/main.ts

# Restart policy
Restart=always
RestartSec=5
StartLimitInterval=60s
StartLimitBurst=3

# Graceful shutdown
TimeoutStopSec=30
KillSignal=SIGTERM

# Logging
StandardOutput=append:%LOGDIR%/app.log
StandardError=append:%LOGDIR%/error.log
SyslogIdentifier=cc-im

# Security - relaxed for user services
NoNewPrivileges=false
ProtectSystem=no
ProtectHome=no

[Install]
# Target is set by install script based on user vs system install
# User install: default.target
# System install: multi-user.target
WantedBy=%TARGET%
